# serdal.service
[Unit]
Description=Serdal.ru Background Workers
After=network.target

[Service]
# Matches your old configuration
User=www-data
Group=www-data

# Path to your project
WorkingDirectory=/var/www/serdal.ru

# Command to run (build + optimize + queue + schedule)
# Ensure /usr/bin/npm is correct. You can check with `which npm`
ExecStart=/usr/bin/npm run prod

# Restart automatically if it crashes
Restart=always
RestartSec=5

# Output logs to syslog
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=serdal-worker

[Install]
WantedBy=multi-user.target
